# trunk/pttbbs/sample/rc.local
#!/bin/sh
# WARNING: this script is run by root!!

/usr/bin/su bbsadm -c /home/bbs/bin/shmctl init
/usr/bin/su bbsadm -c /home/bbs/bin/outmail&
/usr/bin/su bbsadm -c /home/bbs/innd/innbbsd 7777
/home/bbs/bin/bbsctl start
